# resources/serial-keyboard@.service
[Unit]
Description=Serial Keyboard Driver

[Service]
ExecStart=/usr/share/serial-keyboard/start.sh /dev/%i

[Install]
WantedBy=multi-user.target
